# tb2_root_files/etc/udev/rules.d/77-rplidar.rules
# add rplidar with symlink name "/dev/rplidar"
ACTION=="add", ATTRS{idVendor}=="10c4", ATTRS{idProduct}=="ea60", SYMLINK+="rplidar", RUN+="/etc/init.d/rplidar start"
